# deploy/supervisor/dply-log-drain.conf
; dply Logs — app-log drain receiver (UDP syslog → app_logs)
;
; Deploy this ONLY on the box that owns DPLY_LOG_DRAIN_HOST (the control plane,
; or a dedicated drain host). It is a long-lived UDP listener — nothing in the
; web/queue path starts it, so without this program a site's "dply Logs" channel
; sends datagrams into the void and the App logs panel stays empty.
;
; Requires in the app .env on this host (then `php artisan config:cache`):
;   DPLY_LOG_DRAIN_HOST=<this box's address sites send to>
;   DPLY_LOG_DRAIN_PORT=<udp port>
; The port is read from config (log_drains.dply_realtime.port) when --port is omitted.
;
; Firewall: open DPLY_LOG_DRAIN_PORT/udp from customer site servers to this host
; (on-box UFW + provider cloud-firewall). Transport is PLAINTEXT UDP syslog
; routed by the per-site dly_ token — keep the drain host on a trusted/VPC path
; where possible and rely on the firewall to limit senders.
;
; The bash guard keeps the listener supervised without a tight crash-loop when the
; drain port isn't configured yet (it backs off ~30s and retries).

[program:dply-log-drain]
process_name=%(program_name)s
command=/bin/bash -c 'while true; do php %(ENV_DPLY_ROOT)s/artisan dply:log-drain:listen --host=0.0.0.0 || sleep 30; sleep 1; done'
directory=%(ENV_DPLY_ROOT)s
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_DPLY_ROOT)s/storage/logs/log-drain.log
