# deploy/supervisor/dply-worker.conf
; dply control-plane — replica worker (DPLY_RUNTIME=worker, DPLY_WORKER_ROLE=replica)
;
; Env on this host:
;   DPLY_RUNTIME=worker
;   DPLY_WORKER_ROLE=replica
;   HORIZON_NAME=dply-worker-2
;   HORIZON_MAX_PROCESSES=4

[program:dply-horizon]
process_name=%(program_name)s
command=php %(ENV_DPLY_ROOT)s/artisan horizon
directory=%(ENV_DPLY_ROOT)s
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_DPLY_ROOT)s/storage/logs/horizon.log
stopwaitsecs=3660

[program:dply-runtime-check]
process_name=%(program_name)s
command=/bin/bash -c 'while true; do php %(ENV_DPLY_ROOT)s/artisan dply:runtime:check --quiet || logger -t dply-runtime "check failed"; sleep 300; done'
directory=%(ENV_DPLY_ROOT)s
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_DPLY_ROOT)s/storage/logs/runtime-check.log
