# deploy/supervisor/dply-web.conf
; dply control-plane — web tier (DPLY_RUNTIME=web)
;
; Install: copy to /etc/supervisor/conf.d/ and set DPLY_ROOT to your release path.
; Reload: supervisorctl reread && supervisorctl update
;
; Processes: Reverb only. HTTP is nginx + php-fpm (not managed here).
; Horizon and schedule:work must NOT run on this host.

[program:dply-reverb]
process_name=%(program_name)s
command=php %(ENV_DPLY_ROOT)s/artisan reverb:start
directory=%(ENV_DPLY_ROOT)s
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_DPLY_ROOT)s/storage/logs/reverb.log
stopwaitsecs=30
