# deploy/supervisor/dply-worker-primary.conf
; dply control-plane — primary worker (DPLY_RUNTIME=worker, DPLY_WORKER_ROLE=primary)
;
; Env on this host (in addition to shared app .env):
;   DPLY_RUNTIME=worker
;   DPLY_WORKER_ROLE=primary
;   HORIZON_NAME=dply-worker-1
;   HORIZON_MAX_PROCESSES=6
;
; Runs migrations at deploy start before other workers roll (see docs/dply-production-runtime.md).

[program:dply-horizon]
process_name=%(program_name)s
command=php %(ENV_DPLY_ROOT)s/artisan horizon
directory=%(ENV_DPLY_ROOT)s
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_DPLY_ROOT)s/storage/logs/horizon.log
stopwaitsecs=3660

[program:dply-scheduler]
process_name=%(program_name)s
command=php %(ENV_DPLY_ROOT)s/artisan schedule:work
directory=%(ENV_DPLY_ROOT)s
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_DPLY_ROOT)s/storage/logs/scheduler.log
stopwaitsecs=60

[program:dply-runtime-check]
process_name=%(program_name)s
command=/bin/bash -c 'while true; do php %(ENV_DPLY_ROOT)s/artisan dply:runtime:check --quiet || logger -t dply-runtime "check failed"; sleep 300; done'
directory=%(ENV_DPLY_ROOT)s
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_DPLY_ROOT)s/storage/logs/runtime-check.log
